# install/openerp.conf
description "Openerp script de inicializacao"

respawn
respawn limit 15 5

# start in normal runlevels when disks are mounted and networking is available
start on runlevel [2345]

# stop on shutdown/halt, single-user mode and reboot
stop on runlevel [016]

env USER='danimar'
env APP_DIR='/home/danimar/odoo/'
env APP_EXEC='openerp-server'

script

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/openerp.pid --chuid danimar --chdir ${APP_DIR} --exec ${APP_DIR}${APP_EXEC} -- --config=openerp-config

end script
